# make/demo/tms-recycle-demo.crontab
0 7 * * * root sudo -i /usr/sbin/tms-recycle-demo
