# make/packages/tms/lib/systemd/system/tms.service
[Unit]
Description=TMS
After=syslog.target network.target

[Service]
User=prisma
Group=prisma
ExecStart=/usr/sbin/tmsd --start
ExecStop=/usr/sbin/tmsd --stop

[Install]
WantedBy=multi-user.target
